# locales/etc/profile.d/00-i18npath.sh
export I18NPATH='/usr/local/share/i18n/'
